# docker/root/etc/s6-overlay/s6-rc.d/init-bookbagofholding-config/run
#!/command/with-contenv bash
# shellcheck shell=bash

# Create required directories
mkdir -p \
    /config/log \
    /config/cache \
    /downloads \
    /books

# Copy default config if not present
if [[ ! -f /config/config.ini ]]; then
    echo "Creating default config.ini..."
    cp /defaults/config.ini /config/config.ini
fi

# Set permissions
lsiown -R abc:abc /config
lsiown abc:abc /downloads 2>/dev/null || true
lsiown abc:abc /books 2>/dev/null || true

echo "Bookbag of Holding configuration initialized."
